# debian/mpclipboard-server.service
[Unit]
AssertPathExists=/usr/bin/mpclipboard-server

[Service]
WorkingDirectory=~
ExecStart=/usr/bin/mpclipboard-server --start
Restart=always
PrivateTmp=false
NoNewPrivileges=true

[Install]
WantedBy=default.target
